# config/php/supervisord.conf
# =================================================================
# DEVARCH SUPERVISOR CONFIGURATION
# =================================================================
# Supervisor configuration for managing multiple processes in the
# development container including PHP-FPM and development servers.

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700
username=devarch
password=123456

[supervisord]
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
username=devarch
password=123456

# =================================================================
# PHP-FPM PROCESS
# =================================================================
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize --fpm-config /usr/local/etc/php-fpm.conf
user=root
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/php-fpm.log
stderr_logfile=/var/log/supervisor/php-fpm-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
killasgroup=true
stopasgroup=true

# =================================================================
# NODE.JS DEVELOPMENT SERVERS
# =================================================================
# Next.js Development Server (auto-detected)
[program:nextjs-dev]
command=/bin/bash -c 'cd /var/www/html && find . -maxdepth 2 -name "next.config.*" -exec dirname {} \; | head -1 | xargs -I {} bash -c "cd {} && npm run dev"'
user=www-data
autostart=false
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/nextjs-dev.log
stderr_logfile=/var/log/supervisor/nextjs-dev-error.log
environment=NODE_ENV=development,PORT=3000
killasgroup=true
stopasgroup=true

# Generic Node.js Development Server
[program:nodejs-dev]
command=/bin/bash -c 'cd /var/www/html && find . -maxdepth 2 -name "package.json" ! -path "*/node_modules/*" -exec dirname {} \; | head -1 | xargs -I {} bash -c "cd {} && if [ -f package.json ] && grep -q \"dev\" package.json && ! [ -f next.config.js ]; then npm run dev; fi"'
user=www-data
autostart=false
autorestart=true
priority=21
stdout_logfile=/var/log/supervisor/nodejs-dev.log
stderr_logfile=/var/log/supervisor/nodejs-dev-error.log
environment=NODE_ENV=development,PORT=3001
killasgroup=true
stopasgroup=true

# Vite Development Server
[program:vite-dev]
command=/bin/bash -c 'cd /var/www/html && find . -maxdepth 2 -name "vite.config.*" -exec dirname {} \; | head -1 | xargs -I {} bash -c "cd {} && npx vite --host 0.0.0.0 --port 5173"'
user=www-data
autostart=false
autorestart=true
priority=22
stdout_logfile=/var/log/supervisor/vite-dev.log
stderr_logfile=/var/log/supervisor/vite-dev-error.log
environment=NODE_ENV=development
killasgroup=true
stopasgroup=true

# =================================================================
# PYTHON DEVELOPMENT SERVERS
# =================================================================
# Django Development Server (auto-detected)
[program:django-dev]
command=/bin/bash -c 'cd /var/www/html && find . -maxdepth 2 -name "manage.py" -exec dirname {} \; | head -1 | xargs -I {} bash -c "cd {} && python manage.py runserver 0.0.0.0:8000"'
user=www-data
autostart=false
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/django-dev.log
stderr_logfile=/var/log/supervisor/django-dev-error.log
environment=DJANGO_SETTINGS_MODULE=settings,DEBUG=True
killasgroup=true
stopasgroup=true

# Flask Development Server
[program:flask-dev]
command=/bin/bash -c 'cd /var/www/html && find . -maxdepth 2 -name "app.py" -o -name "main.py" | head -1 | xargs -I {} bash -c "cd $(dirname {}) && python {} --host 0.0.0.0 --port 8001"'
user=www-data
autostart=false
autorestart=true
priority=31
stdout_logfile=/var/log/supervisor/flask-dev.log
stderr_logfile=/var/log/supervisor/flask-dev-error.log
environment=FLASK_ENV=development,FLASK_DEBUG=1
killasgroup=true
stopasgroup=true

# FastAPI Development Server (Uvicorn)
[program:fastapi-dev]
command=/bin/bash -c 'cd /var/www/html && find . -maxdepth 2 -name "main.py" -exec dirname {} \; | head -1 | xargs -I {} bash -c "cd {} && uvicorn main:app --host 0.0.0.0 --port 8002 --reload"'
user=www-data
autostart=false
autorestart=true
priority=32
stdout_logfile=/var/log/supervisor/fastapi-dev.log
stderr_logfile=/var/log/supervisor/fastapi-dev-error.log
environment=PYTHONPATH=/var/www/html
killasgroup=true
stopasgroup=true

# =================================================================
# BACKGROUND WORKERS
# =================================================================
# Laravel Queue Worker
[program:laravel-queue]
command=/bin/bash -c 'cd /var/www/html && find . -maxdepth 2 -name "artisan" -exec dirname {} \; | head -1 | xargs -I {} bash -c "cd {} && php artisan queue:work --sleep=3 --tries=3"'
user=www-data
autostart=false
autorestart=true
priority=40
stdout_logfile=/var/log/supervisor/laravel-queue.log
stderr_logfile=/var/log/supervisor/laravel-queue-error.log
numprocs=2
process_name=%(program_name)s_%(process_num)02d
killasgroup=true
stopasgroup=true

# Celery Worker (for Django/Flask)
[program:celery-worker]
command=/bin/bash -c 'cd /var/www/html && find . -maxdepth 2 -name "celery.py" -exec dirname {} \; | head -1 | xargs -I {} bash -c "cd {} && celery -A . worker --loglevel=info"'
user=www-data
autostart=false
autorestart=true
priority=41
stdout_logfile=/var/log/supervisor/celery-worker.log
stderr_logfile=/var/log/supervisor/celery-worker-error.log
killasgroup=true
stopasgroup=true

# Celery Beat (for scheduled tasks)
[program:celery-beat]
command=/bin/bash -c 'cd /var/www/html && find . -maxdepth 2 -name "celery.py" -exec dirname {} \; | head -1 | xargs -I {} bash -c "cd {} && celery -A . beat --loglevel=info"'
user=www-data
autostart=false
autorestart=true
priority=42
stdout_logfile=/var/log/supervisor/celery-beat.log
stderr_logfile=/var/log/supervisor/celery-beat-error.log
killasgroup=true
stopasgroup=true

# =================================================================
# FILE WATCHERS
# =================================================================
# SASS/SCSS Watcher
[program:sass-watcher]
command=/bin/bash -c 'cd /var/www/html && find . -maxdepth 3 -name "*.scss" -o -name "*.sass" | head -1 | xargs -I {} bash -c "cd $(dirname {}) && sass --watch .:. --style expanded"'
user=www-data
autostart=false
autorestart=true
priority=50
stdout_logfile=/var/log/supervisor/sass-watcher.log
stderr_logfile=/var/log/supervisor/sass-watcher-error.log
killasgroup=true
stopasgroup=true

# TypeScript Watcher
[program:typescript-watcher]
command=/bin/bash -c 'cd /var/www/html && find . -maxdepth 2 -name "tsconfig.json" -exec dirname {} \; | head -1 | xargs -I {} bash -c "cd {} && tsc --watch"'
user=www-data
autostart=false
autorestart=true
priority=51
stdout_logfile=/var/log/supervisor/typescript-watcher.log
stderr_logfile=/var/log/supervisor/typescript-watcher-error.log
killasgroup=true
stopasgroup=true

# =================================================================
# MONITORING AND UTILITIES
# =================================================================
# Cron Service
[program:cron]
command=/usr/sbin/cron -f
user=root
autostart=true
autorestart=true
priority=60
stdout_logfile=/var/log/supervisor/cron.log
stderr_logfile=/var/log/supervisor/cron-error.log

# Log Rotation
[program:logrotate]
command=/bin/bash -c 'while true; do logrotate -f /etc/logrotate.conf; sleep 3600; done'
user=root
autostart=true
autorestart=true
priority=61
stdout_logfile=/var/log/supervisor/logrotate.log
stderr_logfile=/var/log/supervisor/logrotate-error.log

# =================================================================
# GROUP DEFINITIONS
# =================================================================
# Development servers group
[group:dev-servers]
programs=nextjs-dev,nodejs-dev,vite-dev,django-dev,flask-dev,fastapi-dev
priority=20

# Background workers group
[group:workers]
programs=laravel-queue,celery-worker,celery-beat
priority=40

# File watchers group
[group:watchers]
programs=sass-watcher,typescript-watcher
priority=50

# System services group
[group:system]
programs=php-fpm,cron,logrotate
priority=10

# =================================================================
# EVENT LISTENERS
# =================================================================
# Auto-restart failed processes
[eventlistener:crashmail]
command=/usr/local/bin/crashmail
events=PROCESS_STATE_EXITED
buffer_size=50
stdout_logfile=/var/log/supervisor/crashmail.log
stderr_logfile=/var/log/supervisor/crashmail-error.log

# Memory monitor
[eventlistener:memmon]
command=/usr/local/bin/memmon -p php-fpm=200MB -m root@localhost
events=TICK_60
buffer_size=10
stdout_logfile=/var/log/supervisor/memmon.log
stderr_logfile=/var/log/supervisor/memmon-error.log

# =================================================================
# INET HTTP SERVER (for web interface)
# =================================================================
[inet_http_server]
port=*:9001
username=devarch
password=123456

# =================================================================
# CONFIGURATION NOTES
# =================================================================
# This supervisor configuration provides:
# 1. Automatic process management for development servers
# 2. Framework detection and auto-starting of appropriate dev servers
# 3. Background worker management for queues and scheduled tasks
# 4. File watchers for SASS/SCSS and TypeScript compilation
# 5. System service management (cron, log rotation)
# 6. Web interface available at port 9001
#
# To control services:
# - Start all: supervisorctl start all
# - Start group: supervisorctl start dev-servers:*
# - Stop service: supervisorctl stop nextjs-dev
# - View status: supervisorctl status
# - Reload config: supervisorctl reread && supervisorctl update
#
# Services are auto-detected based on project structure:
# - Next.js: Looks for next.config.js/mjs
# - Django: Looks for manage.py
# - Laravel: Looks for artisan
# - Generic Node.js: Looks for package.json with "dev" script
# - Vite: Looks for vite.config.*
# - FastAPI: Looks for main.py with FastAPI app